# systemd/racingdsx.service.in
[Unit]
Description=RacingDSX native headless telemetry bridge
Requires=hefesto-dualsense4unix.service
After=hefesto-dualsense4unix.service
PartOf=hefesto-dualsense4unix.service

[Service]
Type=simple
ExecStart="@DOTNET_BIN@" "@REPO_ROOT@/.runtime/racingdsx/RacingDSX.dll"
WorkingDirectory=@REPO_ROOT@/.runtime/racingdsx
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1
Restart=on-failure
RestartSec=2
KillMode=control-group

[Install]
WantedBy=default.target
